# mpd.service.in
[Unit]
Description=Music Player Daemon
After=network.target sound.target

[Service]
ExecStart=@prefix@/bin/mpd --no-daemon

[Install]
WantedBy=multi-user.target
