# systemd/system/mpd.socket.in
[Socket]
@listen_streams@
Backlog=5
KeepAlive=true
PassCredentials=true

[Install]
WantedBy=sockets.target
