# systemd/system/mpd.socket
[Socket]
ListenStream=/run/mpd/socket
ListenStream=6600
Backlog=5
KeepAlive=true
PassCredentials=true

[Install]
WantedBy=sockets.target
